# scripts/ops/templates/cron/ops_daily_report.cron
# zhilian-os daily ops patrol (run as root or deployment user)
# 8:05 every day, append logs to logs/ops/cron.log
5 8 * * * cd /opt/zhilian-os && BASE_URL=http://127.0.0.1 TOKEN=<admin_jwt> bash scripts/ops/daily_ops_report.sh >> logs/ops/cron.log 2>&1
